# redis.supervisord.conf
[program:redis]
command=/usr/bin/redis-server /redis/redis.conf
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
directory=/redis